# roles/common-pre/files/ansible-hourly.cron
#!/bin/bash

ANSIBLE_DIR='/home/chris/.ansible'

ansible-playbook ${ANSIBLE_DIR}/site.yml
