# roles/common/files/ansible-daily.cron
#!/bin/bash

# Cron job to cleanup the Ansible tmp dir

# Not using vars here - broken vars could delete important stuff
# and I've done that before.  I have the fear.

if [[ -d /home/chris/.ansible/tmp/ ]]
then
  tmpwatch 7d /home/chris/.ansible/tmp/
fi
